# admin/deploy/supervisor-neogtb-worker.conf
; ============================================================================
; Supervisor config pour le worker queue NeoGTB Admin
; ----------------------------------------------------------------------------
; Installation sur le VPS :
;   sudo cp deploy/supervisor-neogtb-worker.conf /etc/supervisor/conf.d/neogtb-worker.conf
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start neogtb-worker:*
;
; Vérifier le statut :
;   sudo supervisorctl status neogtb-worker:*
;
; Voir les logs :
;   tail -f /var/www/neogtb-admin/storage/logs/worker.log
;
; ⚠️ Adapter les chemins si besoin :
;   - command : chemin absolu vers PHP + artisan
;   - directory : racine du projet Laravel
;   - user : utilisateur web (souvent www-data)
;   - stdout_logfile : doit être writable par l'utilisateur
; ============================================================================

[program:neogtb-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/neogtb-admin/artisan queue:work --tries=3 --timeout=300 --sleep=3 --max-time=3600
directory=/var/www/neogtb-admin
autostart=true
autorestart=true
startsecs=5
startretries=3
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/neogtb-admin/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=600
stopasgroup=true
killasgroup=true
